# com.sergivalverde.workspace-status-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sergivalverde.workspace-status-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tensor/dev/gh-issues-sync/sync-workspace-status.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/tensor/Library/Logs/workspace-status-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tensor/Library/Logs/workspace-status-sync.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
